# contrib/docker-plugin/systemd/docker-volume-rclone.socket
[Unit]
Description=Docker Volume Plugin for rclone

[Socket]
ListenStream=/run/docker/plugins/rclone.sock

[Install]
WantedBy=sockets.target
