# vendor/github.com/okzk/sdnotify/sample/sample.service
[Unit]
Description=sample

[Service]
Type=notify
ExecStart=/path/to/sample
ExecStop=/bin/kill -SIGTERM $MAINPID
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy = multi-user.target
